# EASY_APM/CONF/rc.local
#!/bin/sh
touch /var/lock/subsys/local
/usr/bin/rdate -s time.bora.net &
export TZ="/etc/localtime"
/usr/local/proftpd/sbin/proftpd
/usr/local/httpd2/bin/apachectl start
/etc/init.d/mysql.server start
